# srv-nfs-101.automount
#############################################
#                                           #
#   Automount: VM data                      #
#                                           #
#   https://github.com/dehesselle/virtctl   #
#                                           #
#############################################
[Unit]
Description=KVM101 data volume

[Automount]
Where=/srv/nfs/101

[Install]
WantedBy=hypervisor.target
